# 51-qubes-app-macchanger.rules
SUBSYSTEM=="net", ACTION=="add", DRIVERS="iwlwifi", RUN+="/usr/local/libexec/qubes/change-mac.sh %k"
